# backend/marketplace-watcher/com.primeport.watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- The watcher runs HERE on the operator's Mac and welds to the Render
     backend: the OKX agent identity (email login, agent #5021) only exists
     in this machine's onchainos session, and API-key logins land in a
     separate empty account (confirmed 2026-07-13), so the cloud cannot hold
     it. caffeinate -s keeps the Mac awake while on AC power so vending
     survives the lid; on battery, sleep pauses vending until wake.

     Install:
       cp com.primeport.watcher.plist ~/Library/LaunchAgents/
       launchctl load ~/Library/LaunchAgents/com.primeport.watcher.plist
     Logs land in marketplace-watcher/data/watcher.{out,err}.log -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.primeport.watcher</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>exec caffeinate -s node watcher.mjs run</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/user/Documents/prime-port/backend/marketplace-watcher</string>
  <key>EnvironmentVariables</key>
  <dict>
    <!-- launchd's login shell skips .zshrc, so fnm's node bin (where okx-a2a
         lives) never reaches PATH on its own; pin the whole PATH here. -->
    <key>PATH</key>
    <string>/Users/user/.local/share/fnm/node-versions/v22.22.2/installation/bin:/Users/user/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <!-- The deployed board: the one existing Render service (free hours are
         capped, so this URL is it until at least August). -->
    <key>BACKEND_URL</key>
    <string>https://prime-port-latest.onrender.com</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/user/Documents/prime-port/backend/marketplace-watcher/data/watcher.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/user/Documents/prime-port/backend/marketplace-watcher/data/watcher.err.log</string>
</dict>
</plist>
